# demo/fs/etc/inittab
::sysinit:/etc/init.d/rcS
::respawn:-login -f guest
::restart:/sbin/init
